# claude-email-agent.service
[Unit]
Description=Claude Email Agent
After=network.target

[Service]
Type=simple
User=timlihk
WorkingDirectory=/home/timlihk/claude-email-agent
ExecStart=/usr/bin/node /home/timlihk/claude-email-agent/claude-email-agent.js
ExecStartPre=/home/timlihk/claude-email-agent/log-rotation.sh
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PATH=/home/timlihk/.npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/timlihk
StandardOutput=append:/home/timlihk/claude-email-agent/agent.log
StandardError=append:/home/timlihk/claude-email-agent/agent-error.log

[Install]
WantedBy=multi-user.target